# thucduong.service
[Unit]
Description=Thuc Duong Backend Spring Boot Application
After=network.target
Wants=network.target

[Service]
Type=forking
User=deploy
Group=deploy
WorkingDirectory=/home/deploy/thucduong
ExecStart=/usr/bin/docker-compose up -d
ExecStop=/usr/bin/docker-compose down
ExecReload=/usr/bin/docker-compose restart
Restart=always
RestartSec=10

# Environment
Environment=SPRING_PROFILES_ACTIVE=prod

# Security
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/home/deploy/thucduong

[Install]
WantedBy=multi-user.target 